# deploy/content-marketing.service
[Unit]
Description=AI Content Marketing App
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/ai-content-marketing
ExecStart=/root/ai-content-marketing/venv/bin/python main.py start
Restart=on-failure
RestartSec=60
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
